# chezmoi/dot_config/autostart/darkman.desktop
[Desktop Entry]
Exec=/run/current-system/sw/bin/darkman run
Name=darkman
OnlyShowIn=i3
Type=Application
