# chezmoi/dot_config/autostart/xsettingsd.desktop
[Desktop Entry]
Exec=/run/current-system/sw/bin/xsettingsd
Name=xsettingsd
OnlyShowIn=i3
Type=Application
